# usr/lib/systemd/system/getty@tty1.service.d/override.conf
[Unit]
Description=Chromium Kiosk Mode
After=network.target

[Service]
User=chromium-kiosk
Group=chromium-kiosk
ExecStartPre=/bin/bash -c 'while ! ping -c1 google.com &>/dev/null; do sleep 1; done'
ExecStart=/usr/bin/chromium-browser --kiosk --allow-running-insecure-content --noerrdialogs --disable-restore-session-state --disable-session-crashed-bubble --disable-infobars --check-for-update-interval=604800 --ignor-gpu-blacklist --disable-gpu disable-gpu-compositin --hide-scrollbars   --incognito --disable-pinch --overscroll-history-navigation=0 --disable-translate --disable-translate-new-ux --disable-new-avatar-menu --no-first-run --fast --fast-start --disable-features=TranslateUI,AutofillServerCommunication
Restart=always
Environment="DISPLAY=:0"

[Install]
WantedBy=multi-user.target
